# ikiwiki/supervisord.conf
[program:spawn-fcgi]
command=spawn-fcgi -s /tmp/fcgi.socket -M 0662 -n -- /usr/sbin/fcgiwrap
user=ikiwiki

[program:nginx]
command=nginx -g "daemon off;"